# upstart/adhd.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Installed by ADHD package.

description     "A/V hardware management"
author          "chromium-os-dev@chromium.org"

# This job maintains information about the A/V hardware which is
# presented to the OS.
#
# It communicates with cras, the audio server, to provide information
# about the devices present.
#

start on started cras
stop on stopping cras
respawn

expect fork
exec /usr/bin/gavd              # Google A/V Daemon
